# systemd/data-to-timesketch.service
[Unit]
Description=Process data, push to Timesketch

[Service]
ExecStart=/bin/bash /opt/watch-to-timesketch.sh

[Install]
WantedBy=multi-user.target
